# infrastructure/my-prime-jennie.service
# /etc/systemd/system/my-prime-jennie.service
[Unit]
Description=my-prime-jennie AI Trading Stack
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
Environment=AUTOSTART_COMPOSE_PROFILES=infra real
WorkingDirectory=/home/youngs75/projects/my-prime-jennie
ExecStart=/home/youngs75/projects/my-prime-jennie/scripts/systemd_autostart.sh
ExecStop=/usr/bin/docker compose --profile real down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
